# config/amslpr.service
[Unit]
Description=VisiGate - Automated License Plate Recognition
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/visigate
ExecStart=/opt/visigate/venv/bin/python -m src.main --web-only
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment="PYTHONPATH=/opt/visigate"
Environment="VisiGate_CONFIG_DIR=/etc/visigate"
Environment="VisiGate_LOG_DIR=/var/log/visigate"
Environment="VisiGate_DATA_DIR=/var/lib/visigate"

[Install]
WantedBy=multi-user.target
