# config/amslpr.service
[Unit]
Description=AMSLPR - Automated License Plate Recognition
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/amslpr
ExecStart=/opt/amslpr/venv/bin/python -m src.main --web-only
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment="PYTHONPATH=/opt/amslpr"
Environment="AMSLPR_CONFIG_DIR=/etc/amslpr"
Environment="AMSLPR_LOG_DIR=/var/log/amslpr"
Environment="AMSLPR_DATA_DIR=/var/lib/amslpr"

[Install]
WantedBy=multi-user.target
